# collectd.conf
[program:collectd]
user=collectd
directory=/
command=collectd -C /etc/collectd/collectd.conf -f
